# ansible/ansible-apps/files/api.service
[Unit]
Description=API server

[Service]
ExecStart=/tmp/api

[Install]
Alias=api.service
